# deploy/thamdinh-api.service
[Unit]
Description=Appraisal Flask API and React SPA
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/app
EnvironmentFile=-/root/app/API.env
ExecStart=/root/app/.venv/bin/gunicorn --workers 2 --timeout 180 --bind 127.0.0.1:5000 --access-logfile - --error-logfile - api.run:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
